# systemctl/media-backup.mount
[Unit]
Description=Additional drive

[Mount]
What=/dev/sda1
Where=/mnt/bkup
Type=exfat
Options=defaults

[Install]
WantedBy=multi-user.target