# systemctl/media-backup.automount
[Unit]
Description=Automount HDD Backup

[Automount]
Where=/mnt/bkup

[Install]
WantedBy=multi-user.target