# systemctl/media-backup.automount
[Unit]
Description=Automount HDD Backup

[Automount]
Where=/media/backup

[Install]
WantedBy=multi-user.target